# _examples/systemd/fail2ban_exporter.service
[Unit]
Description=Fail2Ban Exporter

[Service]
User=fail2ban_exporter
ExecStart=/usr/sbin/fail2ban_exporter

[Install]
WantedBy=multi-user.target